# .config/systemd/user/projects-git-autosync.timer
[Unit]
Description=Run Projects Git auto-sync every 5 minutes

[Timer]
OnBootSec=2min
OnUnitActiveSec=5min
AccuracySec=30s
Persistent=true

[Install]
WantedBy=timers.target
